# runit/services/swarm-login/run
#!/bin/sh

exec 2>&1

cd /opt/swarm/modal-login
yarn dev > /dev/null 2>&1
